# HOME/.config/systemd/user/default.target
[Unit]
Description=Systemd defaults user units
Wants=graphical.service

